# templates/supervisor.ini
[program:flask]
directory=/webapps/devops
command=/usr/bin/python app.py
autostart=true
autorestart=true