# templates/supervisor.ini
[program:app]
command=/usr/bin/uwsgi
  --plugin http,python
  --http :5000
  --logto {{ workdir }}/app.log
  --wsgi-file {{ workdir }}/{{ wsgi_file }}
  --callable app
  --max-requests 1000
  --master
  --processes 1
  --chmod
directory={{ workdir }}
autostart=true
autorestart=true